# kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
autostart=true
autorestart=true
priority=0
auto.create.topics.enable=true
