# kafka/supervisor/topics.conf
[program:create_topics]
directory=/usr/bin/
stopsignal=KILL
autorestart=false
startsecs=0
priority=20
redirect_stderr=true
command=/usr/bin/create_topics.py
